# roles/vault/templates/vault-backup.timer.j2
[Unit]
Description=Vault backup to R2 timer

[Timer]
OnCalendar={{ vault_backup_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
